# trunk/kommunikator/debian/postinst
#!/bin/sh

set -e
#set -x

. /usr/share/debconf/confmodule
. /usr/share/dbconfig-common/dpkg/postinst.mysql 
dbc_generate_include=php:/var/www/kommunikator/config-db.php
dbc_go kommunikator $@

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>

case "$1" in
    configure)


            # Create link to config snippet
            if [ ! -f /etc/apache2/conf.d/kommunikator.conf ] && [ ! -h /etc/apache2/conf.d/kommunikator.conf ]; then
                ln -s /etc/yate-service/apache.conf /etc/apache2/conf.d/kommunikator.conf
            fi

            # Restart the webserver
            if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d apache2 restart >&2
            else
                /etc/init.d/apache2 restart >&2
            fi
        done

        ;;

    abort-upgrade|abort-remove|abort-deconfigure)

        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.



exit 0
