# trunk/kommunikator/debian/prerm
#!/bin/sh

set -e
#set -x

case "$1" in
    remove)

            if [ -e /etc/apache2/conf.d/kommunikator.conf ]; then
                rm /etc/apache2/conf.d/kommunikator.conf
            fi

	    for f in /etc/kommunikator/yate/*
	    do
		ff=/etc/yate/$(basename $f)
		if [ -f $ff ]; then
		    rm $ff
		    dpkg-divert --remove --package $DPKG_MAINTSCRIPT_PACKAGE --rename --divert $ff.old $ff
		fi
	    done

            # Restart daemons affected
            if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d apache2 restart >&2
                invoke-rc.d yate restart >&2
            else
                /etc/init.d/apache2 restart >&2
                /etc/init.d/yate restart >&2
            fi

        ;;

    abort-upgrade|abort-remove|abort-deconfigure)

        ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.



exit 0
