# trunk/kommunikator-devadmin/debian/postinst
#!/bin/sh

set -e
#set -x

case "$1" in
    configure)
	    chgrp www-data /etc/network/interfaces
	    chmod 660 /etc/network/interfaces

            # Create link to config snippet
            if [ ! -f /etc/apache2/conf.d/kommadmin.conf ] && [ ! -h /etc/apache2/conf.d/kommadmin.conf ]; then
                ln -s /etc/kommadmin/apache.conf /etc/apache2/conf.d/kommadmin.conf
            fi

            # Restart the webserver
            if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d apache2 restart >&2
            else
                /etc/init.d/apache2 restart >&2
            fi
        done

        ;;

    abort-upgrade|abort-remove|abort-deconfigure)

        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.



exit 0
